# config/docker-ns/system/nfs-music.automount
[Automount]
Where=/nfs/music
TimeoutIdleSec=600

[Install]
WantedBy=multi-user.target
