# config/docker-ns/system/nfs-transmission.automount
[Automount]
Where=/nfs/transmission
TimeoutIdleSec=600

[Install]
WantedBy=multi-user.target
