# config/docker-ns/system/nfs-music.mount
[Mount]
What=fs.static.lan:/vol/music
Where=/nfs/music
Type=nfs
Options=ro

[Install]
WantedBy=multi-user.target
