# config/docker-ns/system/nft-transmission.mount
[Mount]
What=fs.static.lan:/vol/torrent/transmission
Where=/nfs/transmission
Type=nfs

[Install]
WantedBy=multi-user.target
