# files/config/zfs/50-hdparm.rules
ACTION=="add|change", KERNEL=="sd*", RUN+="/sbin/hdparm -B 127 -S 241 /dev/%k"
